# services/supervisor/conf.d/php-fpm.ini
[program:php-fpm]
command=php-fpm
process_name=%(program_name)s
autostart=true
autorestart=true
startretries=5
exitcodes=0,2,70
stopsignal=QUIT
stopwaitsecs=2
