# services/supervisor/conf.d/queue.ini
[program:email_queue_worker]
command=php /www/influence/think queue:work --queue email
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/email_queue_worker.log
user=root
startsecs=0
stopwaitsecs=60

[program:shopify_queue_worker]
command=php /www/influence/think queue:work --queue shopify_webhooks
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/shopify_queue_worker.log
user=root
startsecs=0
stopwaitsecs=60

[program:shopline_queue_worker]
command=php /www/influence/think queue:work --queue shopline_webhooks
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/shopline_queue_worker.log
user=root
startsecs=0
stopwaitsecs=60

[program:bright_queue_worker]
command=php /www/influence/think queue:work --queue bright_data_webhooks
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/bright_queue_worker.log
user=root
startsecs=0
stopwaitsecs=60
